# ffbg.service
[Unit]
Description=Firefox background to BING

[Service]
Type=simple
ExecStart=%h/.local/bin/bingcopier.sh

[Install]
WantedBy=default.target
